# doc/systemd/go-matter-pack.service
[Unit]
Description=go-matter-pack service
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
# Adjust ExecStart to the command you want to keep running.
# Example: /usr/local/bin/go-matter-pack run
ExecStart=/usr/local/bin/go-matter-pack run
Restart=on-failure
RestartSec=5s
DynamicUser=yes
StateDirectory=go-matter-pack
Environment=XDG_STATE_HOME=/var/lib

[Install]
WantedBy=multi-user.target
